# charts_supervisor.conf
[program:charts]
command = command = </path/to/virtualenv/bin/python> </path/to/virtualenv/bin/gunicorn> charts:app localhost:5000
directory = </path/to/app>
user = <your-username>